# dosthol.service
[Unit]
Description=dosthol (Do something on LAN)

[Service]
ExecStart=/usr/local/bin/dosthold.sh
RestartSec=1
Restart=always

[Install]
WantedBy=multi-user.target
